# deploy/systemd/reqops-agent-sdk.service
[Unit]
Description=agent-sdk orchestration server
After=reqops-postgres.service
Requires=reqops-postgres.service
PartOf=reqops-stack.target

[Service]
Type=simple
EnvironmentFile=-/etc/reqops-stack.env
ExecStartPre=/bin/bash -lc 'cd "${REQOPS_AGENTS_ROOT}/agent-sdk" && test -x .venv/bin/python || (python3 -m venv .venv && .venv/bin/pip install -e .)'
ExecStart=/bin/bash -lc 'cd "${REQOPS_AGENTS_ROOT}/agent-sdk" && export PYTHONPATH=src DATABASE_URL="postgresql://${POSTGRES_USER:-postgres}:${POSTGRES_PASSWORD}@127.0.0.1:${POSTGRES_PORT:-5432}/agent_sdk_server" AGENT_SDK_LOCAL_VOL_ROOT="${AGENT_SDK_LOCAL_VOL_ROOT:-/var/lib/reqops/agent-sdk-volumes}" && exec .venv/bin/python -m uvicorn api.server:app --host 0.0.0.0 --port 7778'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=reqops-stack.target
